# rootfs/etc/service/dnscrypt-proxy/run

#!/bin/sh
exec /bin/dnscrypt-proxy -config /config/dnscrypt-proxy/dnscrypt-proxy.toml
