# rootfs/etc/service/smartdns/run
#!/bin/sh
exec /bin/smartdns.bin -c /config/smartdns/smartdns.conf -f
