# rootfs/etc/service/adguarddns/run
#!/bin/sh
exec /bin/AdGuardHome -c /config/adguardhome/AdGuardHome.yaml -w /config/adguardhome  -h 0.0.0.0 -p 3000
